# unifi.service.d/chown-stuff.conf
[Service]
ExecStartPre=/bin/chown -R unifi /var/lib/unifi
ExecStartPre=/bin/chown -R unifi /var/log/unifi
ExecStartPre=/bin/chown -R unifi /var/run/unifi

