# unifi.service.d/create-missing-dirs.conf
[Service]
ExecStartPre=/bin/bash -c 'cd /usr/lib/unifi && mkdir work && chown -R unifi work'

