# unifi.service.d/custom-mongod.conf
[Service]
ExecStartPre=/bin/bash -c 'cd /usr/lib/unifi/bin && rm ./mongod && ln -s ../../../../opt/mongodb/bin/mongod ./mongod'

